# pkg_scripts/postinst
#!/bin/bash
set -e

systemctl enable board-games-poll.service

exit 0
