# pkg_scripts/prerm
#!/bin/bash
set -e

systemctl disable board-games-poll.service

exit 0
